# deploy/gunicorn.service
[Unit]
Description=Gunicorn daemon for Iffaa Accounting (Django)
After=network.target postgresql.service

[Service]
Type=notify
User=nsl
Group=nsl
WorkingDirectory=/var/www/html/nsl-iffaa-application/backend
Environment="PATH=/var/www/html/nsl-iffaa-application/backend/.venv/bin"
EnvironmentFile=/var/www/html/nsl-iffaa-application/backend/.env

ExecStart=/var/www/html/nsl-iffaa-application/backend/.venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8003 \
    --timeout 60 \
    config.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
